# launchd/com.rocky.dotfiles-backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.rocky.dotfiles-backup</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/rocky/code/dotfiles/scripts/backup-dotfiles.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>EnvironmentVariables</key>
    <dict>
        <key>SSH_AUTH_SOCK</key>
        <string>/Users/rocky/Library/Group Containers/2BUA8C4S2C.com.1password/t/agent.sock</string>
    </dict>
    <key>StandardOutPath</key>
    <string>/tmp/dotfiles-backup.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/dotfiles-backup.log</string>
</dict>
</plist>
